# packer-ansible-yc-flask/ansible/roles/app/templates/app.service.j2
[Unit]
Description=Flask Application
After=network.target

[Service]
Type=simple
User=app-user
WorkingDirectory=/opt/app
Environment=PATH=/opt/app/venv/bin
ExecStart=/opt/app/venv/bin/gunicorn --bind 0.0.0.0:8000 main:app
Restart=always

[Install]
WantedBy=multi-user.target
